# raspios/overlay-rootfs/etc/systemd/system/userconfig.service.d/override.conf
[Service]
# Clear the original ExecStart line
ExecStart=
# Override with a no-op command
ExecStart=/bin/true
